# room/config/sniffer.conf
[program:sniffer]
directory=/home/fi
user=root
command=python3 sniffer.py
autostart=true
autorestart=true